# etc/systemd/system/pystats.service
[Unit]
Description=pystats

[Service]
ExecStart=/usr/local/bin/pystats.py

[Install]
WantedBy=multi-user.target
